# omnibus/config/templates/init-scripts-agent/systemd.checks.service.erb
[Unit]
Description=Datadog Checks Agent
After=network.target datadog-agent.service
BindsTo=datadog-agent.service

[Service]
Type=simple
User=dd-agent
Restart=on-failure
EnvironmentFile=<%= etc_dir %>/environment
ExecStart=/opt/datadog-agent/embedded/bin/checks-agent run --cfgpath=<%= etc_dir %>/datadog.yaml
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
