# q_branch/mcp-evaluation/mcp/mcp-server.service
[Unit]
Description=MCP Evaluation Server
After=network.target

[Service]
Type=simple
User=lima
WorkingDirectory=/mcp
ExecStart=/mcp/bin/mcp-server -mode MODE_PLACEHOLDER -port 8080
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
