# omnibus/config/templates/datadog-puppy/upstart.conf.erb
description "Datadog Agent"

start on started networking
stop on runlevel [!2345]

respawn

setuid dd-agent
console none

script
  exec <%= install_dir %>/bin/agent/agent start -p <%= install_dir %>/run/agent.pid
end script

post-stop script
 rm -f <%= install_dir %>/run/agent.pid
end script
