# omnibus/package-scripts/agent-dmg/postinst
#!/bin/sh
#
# Perform necessary datadog-agent setup steps after package is installed.
#

INSTALL_DIR=/opt/datadog-agent

# macOS-specific variables
OPT_APP_DIR="$INSTALL_DIR/Datadog Agent.app"
APP_DIR="/Applications/Datadog Agent.app"
CONF_DIR=$INSTALL_DIR/etc
RUN_DIR=$INSTALL_DIR/run

# On Mac, the real log folder is located under /opt/datadog-agent/logs, because MacOS upgrades delete /var/log
LOG_DIR=/opt/datadog-agent/logs

# Root-only staging directory shared with the install script.
# Contains env (install-time vars) and saved-conf/ (backup from preinst).
INSTALL_STAGING_DIR="/private/var/root/datadog-install"
INSTALL_ENV_FILE="$INSTALL_STAGING_DIR/env"
SAVED_CONF_DIR="$INSTALL_STAGING_DIR/saved-conf"

# Let's log the standard outputs of this script
LOG_FILE="$LOG_DIR/postinstall.log"
mkdir -vp $LOG_DIR
chmod 750 $LOG_DIR
exec > $LOG_FILE 2>&1

echo "# State at the beginning"
echo "## Agent version"
$INSTALL_DIR/bin/agent/agent version || true
echo "## $INSTALL_DIR"
ls -al $INSTALL_DIR || true
echo "## $APP_DIR/Contents/Resources"
ls -al "$APP_DIR/Contents/Resources" || true

# Stop any running GUI instances before installation/upgrade
# This prevents "Load failed" errors when KeepAlive=true
echo "# Stopping GUI app if running"
for logged_user in $(who | awk '{print $1}' | sort -u); do
    logged_uid=$(id -u "$logged_user" 2>/dev/null) || continue
    launchctl bootout "gui/$logged_uid/com.datadoghq.gui" 2>/dev/null || true
done

# Wait for GUI processes to actually terminate (with 10 second timeout to match ExitTimeOut)
max_wait=100  # 100 * 0.1s = 10 seconds
count=0
while [ $count -lt $max_wait ]; do
    if ! pgrep -f "Datadog Agent.app/Contents/MacOS/gui" > /dev/null 2>&1; then
        echo "GUI app processes terminated"
        break
    fi
    sleep 0.1
    count=$((count + 1))
done
if [ $count -ge $max_wait ]; then
    echo "Warning: GUI processes still running after 10s, proceeding anyway"
fi

echo "# Preparing log dir"
DDAGENT_USER="_dd-agent"

chown -vR "$DDAGENT_USER:admin" "$LOG_DIR"
chmod -v 770 "$LOG_DIR"  # admin group can read and write logs

echo "# Installing the app"
mv -v "$OPT_APP_DIR" /Applications || echo "App already installed"

# Set the run directory for the agent
mkdir -vp "$RUN_DIR"
chown -vR "$DDAGENT_USER:admin" "$RUN_DIR"
chmod -v 775 "$RUN_DIR"  # admin group can read and write runtime files

# Create separate world-writable directory ONLY for multi-user GUI sockets
mkdir -vp "$RUN_DIR/ipc"
# Remove stale GUI socket files before chown to prevent "Address already in use" on reinstall
find "$RUN_DIR/ipc" -maxdepth 1 -name 'gui-*.sock' -exec rm -fv {} \;
chown -vR root:wheel "$RUN_DIR/ipc"
chmod -v 1777 "$RUN_DIR/ipc"  # Sticky bit: users can create GUI sockets but not delete others' files

# Reset location permission prompt attempt count for all users so each user gets a fresh 2 prompt
# attempts after each installation (per-installation behavior). UserDefaults key is in
# ~/Library/Preferences/com.datadoghq.agent.plist and survives reinstall otherwise.
echo "# Resetting location permission prompt attempt count for all users"
for user_home in /Users/*; do
    [ -d "$user_home" ] || continue
    user=$(basename "$user_home")
    [ "$user" = "Shared" ] && continue
    sudo -u "$user" defaults delete com.datadoghq.agent locationPermissionPromptAttemptCount 2>/dev/null || true
done

# Restore config from pre-install (upgrade) or create default (fresh install)
echo "# Setting up configuration"
mkdir -vp $CONF_DIR/checks.d

if [ -d "$SAVED_CONF_DIR" ]; then
    cp -vf "$SAVED_CONF_DIR/datadog.yaml" $CONF_DIR/ 2>/dev/null || true
    cp -vfR "$SAVED_CONF_DIR/conf.d/"* $CONF_DIR/conf.d 2>/dev/null || true
    cp -vn "$SAVED_CONF_DIR/checks.d/"* $CONF_DIR/checks.d 2>/dev/null || true
fi

if [ ! -e "$CONF_DIR/datadog.yaml" ]; then
    sed -E 's/^api_key:$/api_key: APIKEY/' $CONF_DIR/datadog.yaml.example > $CONF_DIR/datadog.yaml
fi

# Read configuration written by the install script (if present).
# When the pkg is installed via GUI Installer or MDM (no install script),
# this file won't exist and the config is left as-is for the admin to configure.
install_method="macos_dmg"
install_method_version="macos_dmg"
if [ -f "$INSTALL_ENV_FILE" ]; then
    dd_api_key=$(grep '^DD_API_KEY=' "$INSTALL_ENV_FILE" | cut -d= -f2-)
    dd_site=$(grep '^DD_SITE=' "$INSTALL_ENV_FILE" | cut -d= -f2-)
    dd_gui_menu=$(grep '^DD_GUI_APP_MENU_ENABLED=' "$INSTALL_ENV_FILE" | cut -d= -f2-)
    dd_install_method=$(grep '^DD_INSTALL_METHOD=' "$INSTALL_ENV_FILE" | cut -d= -f2-)
    dd_install_version=$(grep '^DD_INSTALL_SCRIPT_VERSION=' "$INSTALL_ENV_FILE" | cut -d= -f2-)

    if [ -n "$dd_api_key" ]; then
        sed -i '' "s/^api_key:.*/api_key: $dd_api_key/" "$CONF_DIR/datadog.yaml"
    fi
    if [ -n "$dd_site" ]; then
        sed -i '' "s/^#* *site:.*/site: $dd_site/" "$CONF_DIR/datadog.yaml"
    fi
    if [ -n "$dd_install_method" ]; then
        install_method="$dd_install_method"
        install_method_version="${dd_install_method}-${dd_install_version}"
    fi
fi

install_info_content="---
install_method:
  tool: $install_method
  tool_version: $install_method
  installer_version: $install_method_version
"
echo "$install_info_content" > $CONF_DIR/install_info

echo "# Setting correct rights on conf"
chown -vR "$DDAGENT_USER:admin" $CONF_DIR
chmod -v 770 $CONF_DIR
find $CONF_DIR -type f -exec chmod 660 {} \;

# `datadog-agent` command line
mkdir -vp /usr/local/bin
ln -vs $INSTALL_DIR/bin/agent/agent /usr/local/bin/datadog-agent

if [ ! -e "$CONF_DIR/datadog.yaml" ]; then
    exit 1
fi

# Install agent as system-wide LaunchDaemon
echo "# Installing agent as system-wide LaunchDaemon"
cp -vf "$CONF_DIR/com.datadoghq.agent.plist.example" /Library/LaunchDaemons/com.datadoghq.agent.plist
chown root:wheel /Library/LaunchDaemons/com.datadoghq.agent.plist
chmod 644 /Library/LaunchDaemons/com.datadoghq.agent.plist
launchctl enable system/com.datadoghq.agent 2>/dev/null || true
launchctl bootstrap system /Library/LaunchDaemons/com.datadoghq.agent.plist || true

# Install system-probe as system-wide LaunchDaemon
echo "# Installing system-probe as system-wide LaunchDaemon"
cp -vf "$CONF_DIR/com.datadoghq.sysprobe.plist.example" /Library/LaunchDaemons/com.datadoghq.sysprobe.plist
chown root:wheel /Library/LaunchDaemons/com.datadoghq.sysprobe.plist
chmod 644 /Library/LaunchDaemons/com.datadoghq.sysprobe.plist
launchctl enable system/com.datadoghq.sysprobe 2>/dev/null || true
launchctl bootstrap system /Library/LaunchDaemons/com.datadoghq.sysprobe.plist || true

# Install GUI LaunchAgent for all users (system-wide)
echo "# Configuring GUI app for system-wide installation"
cp -vf "$CONF_DIR/com.datadoghq.gui.plist.example" /Library/LaunchAgents/com.datadoghq.gui.plist

# The GUI plist template includes --headless by default (no menu bar icon).
# DD_GUI_APP_MENU_ENABLED=true removes it to show the menu bar icon.
if [ "$dd_gui_menu" = "true" ]; then
    echo "# Enabling GUI menu bar icon (removing --headless)"
    sed -i '' '/<string>--headless<\/string>/d' /Library/LaunchAgents/com.datadoghq.gui.plist
fi

chown root:wheel /Library/LaunchAgents/com.datadoghq.gui.plist
chmod 644 /Library/LaunchAgents/com.datadoghq.gui.plist

# Load GUI immediately for the console user if they have an active GUI session
CONSOLE_USER=$(stat -f '%Su' /dev/console 2>/dev/null || echo 'root')
CONSOLE_USER_UID=$(id -u "$CONSOLE_USER" 2>/dev/null)
if [ -n "$CONSOLE_USER_UID" ] && sudo -u "$CONSOLE_USER" launchctl managername 2>/dev/null | grep -q "Aqua"; then
    echo "# Loading GUI app for $CONSOLE_USER"
    launchctl enable "gui/$CONSOLE_USER_UID/com.datadoghq.gui" 2>/dev/null || true
    launchctl bootstrap "gui/$CONSOLE_USER_UID" /Library/LaunchAgents/com.datadoghq.gui.plist || true
else
    echo "# No GUI session detected; GUI app will launch at next user login"
fi

echo ""
echo "=========================================="
echo "Datadog Agent GUI (System-Wide)"
echo "=========================================="
echo "The GUI app will launch automatically for each user when they log in."
echo ""
echo "For WiFi data collection (SSID/BSSID), users must grant Location permission:"
echo "  System Settings -> Privacy & Security -> Location Services -> Enable 'Datadog Agent'"
echo ""
echo "Note: The GUI will prompt for this permission when needed."
echo "=========================================="
echo ""

# Clean up the install staging directory (env file with API key, saved config)
rm -rf "$INSTALL_STAGING_DIR"

# A little debriefing won't hurt
echo "# State at the end"
echo "## Agent version"
$INSTALL_DIR/bin/agent/agent version || true
echo "## $INSTALL_DIR"
ls -al $INSTALL_DIR || true
echo "## $APP_DIR/Contents/Resources"
ls -al "$APP_DIR/Contents/Resources" || true

exit 0
