# pkg/fleet/installer/packages/embedded/templates/datadog-agent-ddot.service.tmpl
[Unit]
{{- if .Stable}}
Description=Datadog Distribution of OpenTelemetry Collector
After=datadog-agent.service datadog-agent-exp.service
BindsTo=datadog-agent.service
Conflicts=datadog-agent-exp.service datadog-agent-ddot-exp.service
ConditionPathExists={{.InstallDir}}/embedded/bin/otel-agent
{{- else}}
Description=Datadog Distribution of OpenTelemetry Collector Experiment
After=datadog-agent.service datadog-agent-exp.service
BindsTo=datadog-agent-exp.service
Conflicts=datadog-agent.service datadog-agent-ddot.service
ConditionPathExists={{.InstallDir}}/embedded/bin/otel-agent
{{- end}}

[Service]
Type=simple
PIDFile={{.PIDDir}}/run/otel-agent.pid
User=dd-agent
Restart=on-failure
EnvironmentFile=-{{.EtcDir}}/environment
Environment="DD_FLEET_POLICIES_DIR={{.FleetPoliciesDir}}"
ExecStart={{.InstallDir}}/embedded/bin/otel-agent run --config {{.EtcDir}}/otel-config.yaml --core-config {{.EtcDir}}/datadog.yaml --pidfile {{.PIDDir}}/run/otel-agent.pid
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
