# pkg/fleet/installer/packages/embedded/templates/datadog-agent.service.tmpl
[Unit]
{{- if .Stable}}
Description=Datadog Agent
After=network.target
Conflicts=datadog-agent-exp.service
Wants=datadog-agent-installer.service datadog-agent-trace.service datadog-agent-process.service datadog-agent-sysprobe.service datadog-agent-security.service datadog-agent-ddot.service
{{- else}}
Description=Datadog Agent Experiment
After=network.target
Conflicts=datadog-agent.service
Wants=datadog-agent-installer-exp.service datadog-agent-trace-exp.service datadog-agent-process-exp.service datadog-agent-sysprobe-exp.service datadog-agent-security-exp.service datadog-agent-ddot-exp.service
OnFailure=datadog-agent.service
Before=datadog-agent.service
{{- end}}

[Service]
Type=simple
PIDFile={{.InstallDir}}/run/agent.pid
User=dd-agent
EnvironmentFile=-{{.EtcDir}}/environment
Environment="DD_FLEET_POLICIES_DIR={{.FleetPoliciesDir}}"
{{- if not .Stable}}
Environment="LD_LIBRARY_PATH={{.InstallDir}}/embedded/lib"
{{- end}}
RuntimeDirectory=datadog
StartLimitInterval=10
StartLimitBurst=5
{{- if .Stable}}
Restart=on-failure
ExecStart={{.InstallDir}}/bin/agent/agent run -p {{.InstallDir}}/run/agent.pid
{{- else}}
Restart=no
ExecStart=/usr/bin/timeout --kill-after=15s 3000s {{.InstallDir}}/bin/agent/agent run -p {{.InstallDir}}/run/agent.pid
ExecStopPost=/bin/false
{{- end}}

{{ if .Stable}}
[Install]
WantedBy=multi-user.target
{{- end}}
