# cmd/ai_prompt_logger/com.datadoghq.ai-usage-agent.desktop-monitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.datadoghq.ai-usage-agent.desktop-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/datadog-agent/embedded/bin/ai-usage-agent-native-host</string>
        <string>--desktop-monitor</string>
        <string>--config</string>
        <string>/opt/datadog-agent/etc/ai_usage_native_host.yaml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
</dict>
</plist>
