# omnibus/config/templates/datadog-agent/gui.launchd.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>KeepAlive</key>
        <dict>
            <key>SuccessfulExit</key>
            <false/>
        </dict>
        <key>Label</key>
        <string>com.datadoghq.gui</string>
        <key>ProgramArguments</key>
        <array>
            <string>/Applications/Datadog Agent.app/Contents/MacOS/gui</string>
        </array>
        <key>ExitTimeOut</key>
        <integer>10</integer>
    </dict>
</plist>
