# omnibus/config/templates/datadog-agent/systemd.security-reload.path.erb
[Path]
PathChanged=<%= etc_dir %>/enabled.d/security-agent

[Install]
WantedBy=multi.user.target
