# omnibus/config/templates/datadog-agent/systemd.security.service.erb
[Unit]
Description=Datadog Security Agent
After=network.target datadog-agent.service
BindsTo=datadog-agent.service
ConditionPathExists=<%= etc_dir %>/security-agent.yaml

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/security-agent.pid
Restart=on-failure
ExecStart=<%= install_dir %>/embedded/bin/security-agent start -c <%= etc_dir %>/datadog.yaml -c <%= etc_dir %>/security-agent.yaml -p <%= install_dir %>/run/security-agent.pid
Environment=GOTRACEBACK=crash
# Since systemd 229, should be in [Unit] but in order to support systemd <229,
# it is also supported to have it here.
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
