# omnibus/config/templates/datadog-agent/systemd.service.erb
[Unit]
Description="Datadog Agent"
After=network.target

[Service]
User=root
ExecStart=<%= install_dir %>/bin/agent/agent

[Install]
WantedBy=multi-user.target
