# omnibus/config/templates/datadog-agent/systemd.trace.service.erb
[Unit]
Description="Datadog Trace Agent (APM)"
After=network.target datadog-agent.service
BindsTo=datadog-agent.service
StartLimitInterval=10
StartLimitBurst=5

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/trace-agent.pid
User=dd-agent
Restart=on-failure
ExecStart=<%= install_dir %>/embedded/bin/trace-agent --config <%= etc_dir %>/datadog.yaml --pid <%= install_dir %>/run/trace-agent.pid

[Install]
WantedBy=multi-user.target
