# omnibus/config/templates/datadog-agent/systemd.trace.service.erb
[Unit]
Description=Datadog Trace Agent (APM)
After=network.target datadog-agent.service
BindsTo=datadog-agent.service

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/trace-agent.pid
User=dd-agent
Restart=on-failure
ExecStart=<%= install_dir %>/embedded/bin/trace-agent --config <%= etc_dir %>/datadog.yaml --pid <%= install_dir %>/run/trace-agent.pid
Environment=GOTRACEBACK=crash
# Since systemd 229, should be in [Unit] but in order to support systemd <229,
# it is also supported to have it here.
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
