# omnibus/config/templates/datadog-agent/upstart.conf.erb
description "Datadog Agent"

start on started networking
stop on runlevel [!2345]

respawn

console none

script
  exec <%= install_dir %>/bin/agent/agent start -p /var/run/datadog/agent.pid
end script

post-stop script
 rm -f /var/run/datadog/agent.pid
 rm -f /tmp/agent.sock
end script
