# omnibus/config/templates/datadog-agent/upstart.conf.erb
description "Datadog Agent"

start on started networking
stop on runlevel [!2345]

respawn

console log

script
  # Log to console is disabled since the agent already logs using file or syslog depending on the configuration.
  #
  # setuid is not available in versions of upstart before 1.4. CentOS/RHEL6 use an earlier version of upstart.
  # This is the best way to set the user in the absence of setuid.
  exec su -s /bin/sh -c 'DD_LOG_TO_CONSOLE=false exec "$0" "$@"' dd-agent -- <%= install_dir %>/bin/agent/agent start -p <%= install_dir %>/run/agent.pid
end script

post-stop script
 rm -f <%= install_dir %>/run/agent.pid
 rm -f /tmp/agent.sock
end script
