# omnibus/package-scripts/dogstatsd-rpm/prerm
#!/bin/sh
#
# Perform necessary datadog-dogstatsd setup steps prior to remove the old package.
#
# .rpm: STEP 4 of 6

INSTALL_DIR=/opt/datadog-dogstatsd
SERVICE_NAME=datadog-dogstatsd

stop_agent()
{
    # Stop an already running dogstatsd
    # Only supports systemd and upstart
    if command -v systemctl >/dev/null 2>&1; then
        systemctl stop $SERVICE_NAME || true
    elif command -v initctl >/dev/null 2>&1; then
        initctl stop $SERVICE_NAME || true
    else
        echo "[ WARNING ]\tCannot detect a supported init system. The datadog-dogstatsd package only provides service files for systemd and upstart."
    fi
}

deregister_agent()
{
    # Disable dogstatsd start on system boot
    # Only supports systemd and upstart
    if command -v systemctl >/dev/null 2>&1; then
        systemctl disable $SERVICE_NAME || true
    elif command -v initctl >/dev/null 2>&1; then
        # Nothing to do, this is defined directly in the upstart job file
        :
    else
        echo "[ WARNING ]\tCannot detect a supported init system. The datadog-dogstatsd package only provides service files for systemd and upstart."
    fi
}

cleanup_run_directory()
{
    # Clean up runtime-created files (sockets, pid files, etc.)
    # These are created by dogstatsd at runtime, not by the package.
    rm -f "$INSTALL_DIR/run/dsd.socket" || true
    rm -f "$INSTALL_DIR/run/dogstatsd.pid" || true

    # Remove the run directory if empty
    rmdir "$INSTALL_DIR/run" 2>/dev/null || true
}

stop_agent
deregister_agent

case "$*" in
    0)
        # We're uninstalling.
        cleanup_run_directory
    ;;
    1)
        # We're upgrading.
    ;;
    *)
    ;;
esac

exit 0
