# omnibus/config/templates/updater/datadog-updater.service.erb
[Unit]
Description=Datadog Updater
After=network.target

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/updater.pid
User=dd-updater
Restart=on-failure
EnvironmentFile=-<%= etc_dir %>/environment
ExecStart=<%= install_dir %>/bin/updater/updater run -p <%= install_dir %>/run/updater.pid
# Since systemd 229, should be in [Unit] but in order to support systemd <229,
# it is also supported to have it here.
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
