# omnibus/config/templates/updater/datadog-agent-sysprobe-exp.service.erb
[Unit]
Description=Datadog System Probe Experiment
Requires=sys-kernel-debug.mount
Before=datadog-agent.service
After=network.target sys-kernel-debug.mount
BindsTo=datadog-agent-exp.service
ConditionPathExists=<%= etc_dir %>/system-probe.yaml

[Service]
Type=simple
PIDFile=<%= agent_dir %>/run/system-probe.pid
Restart=on-failure
ExecStart=<%= agent_dir %>/embedded/bin/system-probe run --config=<%= etc_dir %>/system-probe.yaml --pid=<%= agent_dir %>/run/system-probe.pid
# Since systemd 229, should be in [Unit] but in order to support systemd <229,
# it is also supported to have it here.
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
