# omnibus/config/templates/updater/datadog-updater-admin.path.erb
[Unit]
Description="Monitor requests to execute systemd commands"

[Path]
PathModified=<%= install_dir %>/run/in.fifo

[Install]
WantedBy=multi-user.target
