# omnibus/config/templates/updater/datadog-agent.path.erb
[Unit]
Description="Monitor requests to stop experiment"

[Path]
PathModified=<%= install_dir %>/systemd_commands/stop_experiment

[Install]
WantedBy=multi-user.target
