# omnibus/config/templates/updater/datadog-agent-process-exp.service.erb
[Unit]
Description=Datadog Process Agent Experiment
After=network.target
BindsTo=datadog-agent-exp.service

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/process-agent.pid
User=dd-agent
Restart=on-failure
EnvironmentFile=-<%= etc_dir %>/environment
ExecStart=<%= install_dir %>/agent_entrypoints/experiment_agent/embedded/bin/process-agent --cfgpath=<%= etc_dir %>/datadog.yaml --sysprobe-config=<%= etc_dir %>/system-probe.yaml --pid=<%= install_dir %>/run/process-agent.pid
# Since systemd 229, should be in [Unit] but in order to support systemd <229,
# it is also supported to have it here.
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
