# omnibus/config/templates/updater/datadog-agent-exp.service.erb
[Unit]
Description=Datadog Agent Experiment
After=network.target
OnFailure=datadog-agent.service 
Conflicts=datadog-agent.service 
Before=datadog-agent.service 
JobTimeoutSec=3000
Wants=datadog-agent-trace-exp.service datadog-agent-process-exp.service datadog-agent-sysprobe-exp.service datadog-agent-security-exp.service

[Service]
Type=oneshot
PIDFile=<%= agent_dir %>/run/agent.pid
User=dd-agent
EnvironmentFile=-<%= etc_dir %>/environment
ExecStart=<%= agent_dir %>/bin/agent/agent run -p <%= agent_dir %>/run/agent.pid
ExecStart=<%= agent_dir %>/bin/agent/agent run -p <%= agent_dir %>/run/agent.pid
ExecStart=<%= agent_dir %>/bin/agent/agent run -p <%= agent_dir %>/run/agent.pid
ExecStart=/bin/false
ExecStop=/bin/false

[Install]
WantedBy=multi-user.target
