# omnibus/config/templates/updater/datadog-updater-admin.service.erb
[Unit]
Description=Clean start of datadog admin executor

[Service]
Type=oneshot
User=root
ExecStart=systemctl kill -s SIGKILL datadog-updater-admin-exec.service
ExecStart=systemctl start datadog-updater-admin-exec.service

[Install]
WantedBy=multi-user.target
