# omnibus/config/templates/updater/datadog-updater-admin.service.erb
[Unit]
Description=Clean start of datadog admin executor

[Service]
Type=oneshot
User=root
ExecStart=-systemctl kill -s SIGKILL datadog-updater-admin-exec.service
ExecStart=test -e <%= install_dir %>/run/in.fifo
ExecStart=systemctl start datadog-updater-admin-exec.service --no-block

[Install]
WantedBy=multi-user.target
