# omnibus/config/templates/updater/datadog-agent-trace-exp.service.erb
[Unit]
Description=Datadog Trace Agent Experiment
BindsTo=datadog-agent-exp.service

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/trace-agent.pid
User=dd-agent
Restart=on-failure
EnvironmentFile=-<%= etc_dir %>/environment
ExecStart=<%= install_dir %>/agent_entrypoints/experiment_agent/embedded/bin/trace-agent --config <%= etc_dir %>/datadog.yaml --pidfile <%= install_dir %>/run/trace-agent.pid
# Since systemd 229, should be in [Unit] but in order to support systemd <229,
# it is also supported to have it here.
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
