# omnibus/config/templates/updater/datadog-installer-exp.service.erb
[Unit]
Description=Datadog Updater Experiment
After=network.target
OnFailure=datadog-updater.service
Conflicts=datadog-updater.service
JobTimeoutSec=3000 #50 minutes

[Service]
Type=oneshot
PIDFile=<%= install_dir %>/run/updater.pid
User=dd-agent
EnvironmentFile=-<%= etc_dir %>/environment
ExecStart=<%= install_dir %>/bin/updater/updater run -p <%= install_dir %>/run/agent.pid
ExecStart=<%= install_dir %>/bin/updater/updater run -p <%= install_dir %>/run/agent.pid
ExecStart=<%= install_dir %>/bin/updater/updater run -p <%= install_dir %>/run/agent.pid
ExecStart=/bin/false
ExecStop=/bin/false

[Install]
WantedBy=multi-user.target
