# omnibus/config/templates/datadog-dogstatsd/systemd-data-plane.service.erb
[Unit]
Description=Datadog Agent Data Plane (dogstatsd)
After=network.target datadog-dogstatsd.service
BindsTo=datadog-dogstatsd.service

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/agent-data-plane.pid
User=dd-agent
Restart=on-failure
ExecStart=<%= install_dir %>/embedded/bin/agent-data-plane --config /etc/datadog-dogstatsd/dogstatsd.yaml run --pidfile <%= install_dir %>/run/agent-data-plane.pid
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
