# omnibus/config/templates/datadog-dogstatsd/systemd.service.erb
[Unit]
Description=Datadog Statsd Server
After=network.target
Wants=datadog-dogstatsd-data-plane.service

[Service]
Type=simple
User=dd-agent
Restart=on-failure
ExecStart=<%= install_dir %>/bin/dogstatsd start -c /etc/datadog-dogstatsd

[Install]
WantedBy=multi-user.target
