# omnibus/config/templates/datadog-dogstatsd/systemd.service.erb
[Unit]
Description="Stackstate statsd server"
After=network.target

[Service]
Type=simple
User=sts-agent
Restart=on-failure
ExecStart=<%= install_dir %>/bin/dogstatsd start -c /etc/stackstate-dogstatsd

[Install]
WantedBy=multi-user.target
