# test/new-e2e/scenarios/netpath-vm/assets/conn-gen.service
[Unit]
Description=Periodic outbound TCP/UDP probes for NPM + Network Path dynamic tests
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/conn-gen.sh
