# pkg/fleet/installer/packages/embedded/tmpl/gen/debrpm-nocap/datadog-agent-sysprobe.service
[Unit]
Description=Datadog System Probe
Requires=sys-kernel-debug.mount
Before=datadog-agent.service datadog-agent-exp.service
After=network.target sys-kernel-debug.mount
BindsTo=datadog-agent.service
Conflicts=datadog-agent-exp.service datadog-agent-sysprobe-exp.service
ConditionPathExists=|/etc/datadog-agent/system-probe.yaml
ConditionPathExists=|/etc/datadog-agent/managed/datadog-agent/stable/system-probe.yaml

[Service]
Type=simple
PIDFile=/opt/datadog-agent/run/system-probe.pid
Restart=on-failure
EnvironmentFile=-/etc/datadog-agent/environment
Environment="DD_FLEET_POLICIES_DIR=/etc/datadog-agent/managed/datadog-agent/stable"
ExecStart=/opt/datadog-agent/embedded/bin/sd-agent -- /opt/datadog-agent/embedded/bin/system-probe run --config=/etc/datadog-agent/system-probe.yaml --pid=/opt/datadog-agent/run/system-probe.pid
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
