# pkg/fleet/installer/packages/embedded/tmpl/gen/debrpm-nocap/datadog-agent-exp.service
[Unit]
Description=Datadog Agent Experiment
After=network.target
Conflicts=datadog-agent.service
Wants=datadog-agent-installer-exp.service datadog-agent-trace-exp.service datadog-agent-process-exp.service datadog-agent-sysprobe-exp.service datadog-agent-security-exp.service datadog-agent-data-plane-exp.service datadog-agent-ddot-exp.service datadog-agent-action-exp.service datadog-agent-procmgrd-exp.service
OnFailure=datadog-agent.service
Before=datadog-agent.service

[Service]
Type=simple
PIDFile=/opt/datadog-agent/run/agent.pid
User=dd-agent
EnvironmentFile=-/etc/datadog-agent-exp/environment
Environment="DD_FLEET_POLICIES_DIR=/etc/datadog-agent-exp/managed/datadog-agent/stable"
Environment="DD_INVENTORIES_FIRST_RUN_DELAY=5"
Environment="EXPERIMENT_TIMEOUT=3000s"
RuntimeDirectory=datadog
StartLimitInterval=10
StartLimitBurst=5
Restart=no
ExecStart=/usr/bin/timeout --kill-after=15s $EXPERIMENT_TIMEOUT /opt/datadog-agent/bin/agent/agent run -c /etc/datadog-agent-exp --sysprobecfgpath=/etc/datadog-agent-exp/system-probe.yaml -p /opt/datadog-agent/run/agent.pid
ExecStopPost=/bin/false


