# pkg/fleet/installer/packages/embedded/tmpl/gen/oci-nocap/datadog-agent-ddot.service
[Unit]
Description=Datadog Distribution of OpenTelemetry Collector
After=datadog-agent.service datadog-agent-exp.service
BindsTo=datadog-agent.service
Conflicts=datadog-agent-exp.service datadog-agent-ddot-exp.service
ConditionPathExists=/opt/datadog-packages/datadog-agent/stable/ext/ddot/embedded/bin/otel-agent
ConditionPathExists=!/etc/datadog-agent/processes.d/datadog-agent-ddot.yaml

[Service]
Type=simple
PIDFile=/opt/datadog-packages/datadog-agent/stable/run/otel-agent.pid
User=dd-agent
Restart=on-failure
RestartSec=2
EnvironmentFile=-/etc/datadog-agent/environment
Environment="DD_FLEET_POLICIES_DIR=/etc/datadog-agent/managed/datadog-agent/stable"
ExecStart=/opt/datadog-packages/datadog-agent/stable/ext/ddot/embedded/bin/otel-agent run --config /etc/datadog-agent/otel-config.yaml --core-config /etc/datadog-agent/datadog.yaml --pidfile /opt/datadog-packages/datadog-agent/stable/run/otel-agent.pid
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
