# pkg/fleet/installer/packages/embedded/tmpl/gen/oci-nocap/datadog-agent-procmgrd-exp.service
[Unit]
Description=Datadog Process Manager Daemon Experiment
After=datadog-agent.service datadog-agent-exp.service
BindsTo=datadog-agent-exp.service
Conflicts=datadog-agent.service datadog-agent-procmgrd.service
ConditionPathExists=/opt/datadog-packages/datadog-agent/experiment/embedded/bin/dd-procmgrd

[Service]
Type=simple
User=dd-agent
Restart=on-failure
EnvironmentFile=-/etc/datadog-agent-exp/environment
Environment="DD_PM_CONFIG_DIR=/etc/datadog-agent-exp/processes.d"
ExecStart=/opt/datadog-packages/datadog-agent/experiment/embedded/bin/dd-procmgrd
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
