# pkg/fleet/installer/packages/embedded/tmpl/gen/oci-nocap/datadog-agent-rshell-privileged-exp.service
[Unit]
Description=Datadog rshell privileged helper experiment
After=datadog-agent-exp.service datadog-agent-action-exp.service
BindsTo=datadog-agent-exp.service datadog-agent-action-exp.service
Conflicts=datadog-agent-rshell-privileged.service
ConditionPathExists=/etc/datadog-agent-exp/rshell-privileged-helper.json

[Service]
Type=simple
Group=dd-agent
UMask=0077
LoadCredential=rshell-verification.json:/etc/datadog-agent-exp/rshell-privileged-helper.json
ExecStart=/opt/datadog-packages/datadog-agent/experiment/embedded/bin/rshell privileged-helper --credential=%d/rshell-verification.json --idle-timeout=30s
PrivateNetwork=yes
RestrictAddressFamilies=AF_UNIX
LockPersonality=yes
RestrictRealtime=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
