# pkg/fleet/installer/packages/embedded/tmpl/datadog-agent-data-plane.service.tmpl
[Unit]
{{- if .Stable}}
Description=Datadog Agent Data Plane
After=datadog-agent.service datadog-agent-exp.service
BindsTo=datadog-agent.service
Conflicts=datadog-agent-exp.service datadog-agent-data-plane-exp.service
{{- else}}
Description=Datadog Agent Data Plane Experiment
After=datadog-agent.service datadog-agent-exp.service
BindsTo=datadog-agent-exp.service
Conflicts=datadog-agent.service datadog-agent-data-plane.service
{{- end}}

[Service]
Type=simple
PIDFile={{.InstallDir}}/run/agent-data-plane.pid
User=dd-agent
Restart=on-failure
EnvironmentFile=-{{.EtcDir}}/environment
Environment="DD_FLEET_POLICIES_DIR={{.FleetPoliciesDir}}"
ExecStart={{.InstallDir}}/embedded/bin/agent-data-plane --config {{.EtcDir}}/datadog.yaml run --pidfile {{.InstallDir}}/run/agent-data-plane.pid
StartLimitInterval=10
StartLimitBurst=5
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
