# omnibus/config/templates/datadog-cluster-agent/systemd.service.erb
[Unit]
Description="Datadog Cluster Agent"
After=network.target

[Service]
Type=simple
User=dd-agent
Restart=on-failure
ExecStart=<%= install_dir %>/bin/datadog-cluster-agent start -c /etc/datadog-cluster-agent

[Install]
WantedBy=multi-user.target
