# omnibus/config/templates/datadog-iot-agent/systemd.service.erb
[Unit]
Description=Datadog IoT Agent
After=network.target

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/agent.pid
User=dd-agent
Restart=on-failure
ExecStart=<%= install_dir %>/bin/agent/agent run -p <%= install_dir %>/run/agent.pid

[Install]
WantedBy=multi-user.target
