# pkg/fleet/installer/packages/embedded/tmpl/gen/oci/datadog-agent-installer.service
[Unit]
Description=Datadog Installer
After=network.target datadog-agent.service datadog-agent-exp.service
BindsTo=datadog-agent.service
Conflicts=datadog-agent-installer-exp.service
ConditionPathExists=!/etc/systemd/system/datadog-installer.service
ConditionPathExists=!/etc/systemd/system/datadog-installer-exp.service

[Service]
Type=simple
PIDFile=/opt/datadog-packages/datadog-agent/stable/run/installer.pid
Restart=on-failure
EnvironmentFile=-/etc/datadog-agent/environment
Environment="DD_FLEET_POLICIES_DIR=/etc/datadog-agent/managed/datadog-agent/stable"
ExecStart=/opt/datadog-packages/datadog-agent/stable/embedded/bin/installer run -c /etc/datadog-agent -p /opt/datadog-packages/datadog-agent/stable/run/installer.pid
StartLimitInterval=10
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
