# pkg/fleet/installer/packages/embedded/tmpl/gen/oci/datadog-agent.service
[Unit]
Description=Datadog Agent
After=network.target
Conflicts=datadog-agent-exp.service
Wants=datadog-agent-installer.service datadog-agent-trace.service datadog-agent-process.service datadog-agent-sysprobe.service datadog-agent-security.service datadog-agent-data-plane.service datadog-agent-ddot.service datadog-agent-action.service datadog-agent-procmgrd.service

[Service]
Type=simple
PIDFile=/opt/datadog-packages/datadog-agent/stable/run/agent.pid
User=dd-agent
EnvironmentFile=-/etc/datadog-agent/environment
Environment="DD_FLEET_POLICIES_DIR=/etc/datadog-agent/managed/datadog-agent/stable"
RuntimeDirectory=datadog
StartLimitInterval=10
StartLimitBurst=5
AmbientCapabilities=CAP_NET_BIND_SERVICE
Restart=on-failure
ExecStart=/opt/datadog-packages/datadog-agent/stable/bin/agent/agent run -p /opt/datadog-packages/datadog-agent/stable/run/agent.pid


[Install]
WantedBy=multi-user.target
