# pkg/fleet/installer/packages/embedded/tmpl/gen/oci/datadog-agent-rshell-privileged-exp.socket
[Unit]
Description=Datadog rshell privileged helper experiment socket
After=datadog-agent-exp.service datadog-agent-action-exp.service
BindsTo=datadog-agent-exp.service datadog-agent-action-exp.service
Conflicts=datadog-agent-rshell-privileged.socket
ConditionPathExists=/etc/datadog-agent-exp/rshell-privileged-helper.json

[Socket]
ListenStream=/run/datadog/rshell-privileged-exp.sock
Service=datadog-agent-rshell-privileged-exp.service
SocketUser=dd-agent
SocketGroup=dd-agent
SocketMode=0600
DirectoryMode=0750
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
