# pkg/fleet/installer/packages/embedded/tmpl/gen/oci/datadog-agent-action-exp.service
[Unit]
Description=Datadog Private Action Runner Experiment
After=network.target datadog-agent.service datadog-agent-exp.service
BindsTo=datadog-agent-exp.service
Conflicts=datadog-agent.service datadog-agent-action.service

[Service]
Type=simple
User=dd-agent
Restart=on-failure
EnvironmentFile=-/etc/datadog-agent-exp/environment
Environment="DD_FLEET_POLICIES_DIR=/etc/datadog-agent-exp/managed/datadog-agent/stable"
Environment="DD_INVENTORIES_FIRST_RUN_DELAY=5"
ExecStart=/opt/datadog-packages/datadog-agent/experiment/embedded/bin/privateactionrunner run --cfgpath=/etc/datadog-agent-exp/datadog.yaml
StartLimitInterval=10
StartLimitBurst=5
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
